# nfs_exporter.service
[Unit]
Description=Test NFS mounts

[Service]
ExecStart=/usr/local/bin/nfs_exporter -o /var/run/prometheus/metrics/nfs_exporter.prom
